# linux/pi-agent/hearth-pi-agent.service
[Unit]
Description=Hearth-Connect Pi Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/hearth-pi-agent
EnvironmentFile=/opt/hearth-pi-agent/config.env
Environment=CONFIG_FILE=/opt/hearth-pi-agent/config.env
ExecStart=/usr/bin/python3 /opt/hearth-pi-agent/pi-agent.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
